# mezzanine-twitter.cron
*/10 * * * * python /app/manage.py poll_twitter
